# crontab
# Edit this to suit you then Put it in /etc/cron.d/momframe

0 2 * * * lukec /usr/bin/perl /home/lukec/src/momframe/email-pics > /home/lukec/src/momframe/momframe.log 2>&1
